# crontab
TZ=Europe/Madrid
@reboot root cd /app && /usr/bin/python -u checker.py >> /proc/1/fd/1 2>&1
0 8,11,14,17,19 * * * root cd /app && /usr/bin/python -u checker.py >> /proc/1/fd/1 2>&1
